# deployment/systemd/smart-blood-queue-matching.service
[Unit]
Description=Smart Blood Queue Worker (matching)
After=network.target redis.service

[Service]
User=www-data
Group=www-data
Restart=always
RestartSec=5
WorkingDirectory=/var/www/smart-blood
ExecStart=/usr/bin/php artisan queue:work --queue=matching --sleep=1 --tries=5 --backoff=10 --timeout=120 --max-time=3600
ExecReload=/bin/kill -USR2 $MAINPID
KillSignal=SIGTERM
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
